# systemd/openclaw-miniapp-bridge.service
[Unit]
Description=OpenClaw Telegram Mini App Bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/openclaw-telegram-miniapp
ExecStart=/usr/bin/env python3 %h/openclaw-telegram-miniapp/bridge/openclaw_miniapp_bridge.py
Environment="MINIAPP_HOST=127.0.0.1"
Environment="MINIAPP_PORT=8765"
Environment="OPENCLAW_BASE_URL=http://127.0.0.1:18789"
Environment="TELEGRAM_ACCOUNT_NAME=default"
Environment="MINIAPP_INITDATA_MAX_AGE_SECONDS=86400"
# Optional: set TELEGRAM_OWNER_ID to your numeric Telegram user id.
# Optional: set MINIAPP_PUBLIC_ORIGIN to your fixed HTTPS origin for stricter CORS.
# Optional: set MINIAPP_SHARED_TOKEN only if you want a separate browser fallback token.
# Optional: set MINIAPP_AUTH_DEBUG=true only while diagnosing auth.
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
